# supervisor/armada-bind.conf
[program:armada-bind]
directory=/opt/armada-bind/src/
command=python armada-bind.py
